# deploy/avira-support-bot.service
[Unit]
Description=Avira support Telegram bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/avira-bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/avira-bot/venv/bin/python -m src.support_bot
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
